# jobbot.service
[Unit]
Description=AI Job Discovery Bot (Scheduler + Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/Anti-gravity
ExecStart=/home/ec2-user/Anti-gravity/venv/bin/python run_service.py
Restart=always
RestartSec=30
StandardOutput=append:/home/ec2-user/Anti-gravity/.tmp/service.log
StandardError=append:/home/ec2-user/Anti-gravity/.tmp/service.log
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=

[Install]
WantedBy=multi-user.target
